# fauxpenstack.socket
[Socket]
ListenStream=${SVC_PORT}
Accept=no
Service=fauxpenstack.service

[Install]
WantedBy=default.target
