# fauxpenstack.service
[Service]
WorkingDirectory=${PWD}
ExecStart=${PWD}/.venv/bin/fauxpenstack --idle 1800
